# etc/crontab/mirror.cron
# -*- cron -*-
#

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# backup ent every day at 3am 
#min   hour  dom  month dow   command
#0-59  0-23  1-31  1-12 0-7
#
    0     *     *     *   *    cd dv/feynman/etc/crontab && ./mirror.sh

# end of file
