# etc/crontab/mirror.sh
#!/bin/bash
#
# michael a.g. aivazis
# orthologue
# (c) 1998-2014 all rights reserved
#

# bzr
bzr="/usr/bin/bzr"

# pull config
echo pulling
${bzr} pull
echo done

# end of file
